# supervisor/rabbitmq.conf
[program:rabbitmq]
command=/bin/sh /opt/rabbitmq.sh
user=root
autostart=true
autorestart=true
environment=
  HOME="/var/lib/rabbitmq/",
  RABBITMQ_PID_FILE="/var/run/rabbitmq/pid"